# chef/cookbooks/carli-select/templates/carli-select.conf.erb
#!upstart

# TODO install this with Chef
# /etc/init/carli-select.conf

description "CARLI Select"

start on startup
stop on shutdown

expect fork

env NODE_BIN_DIR="/usr/bin"
env NODE_PATH="/usr/lib/node_modules"
env APPLICATION_PATH="/var/www/carli-select/latest/middleware/middleware/index.js"
env PIDFILE="/var/run/carli-select.pid"
env LOG="/var/log/carli-select.log"
env MIN_UPTIME="5000"
env SPIN_SLEEP_TIME="2000"

script
    PATH=$NODE_BIN_DIR:$PATH

    exec forever \
        --pidFilev $PIDFILE \
        -a \
        -l $LOG \
        --minUptime $MIN_UPTIME \
        --spinSleepTime $SPIN_SLEEP_TIME \
        start $APPLICATION_PATH
end script

pre-stop script
    PATH=$NODE_BIN_DIR:$PATH
    exec forever stop $APPLICATION_PATH
end script
